# launchd/com.kimgwonil.youtube-shorts-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kimgwonil.youtube-shorts-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/kimgwonil/youtube-shorts-automation &amp;&amp; source .venv/bin/activate &amp;&amp; python scripts/run_daily.py &gt;&gt; logs/daily.log 2&gt;&amp;1</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kimgwonil/youtube-shorts-automation</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kimgwonil/youtube-shorts-automation/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kimgwonil/youtube-shorts-automation/logs/launchd.err.log</string>
</dict>
</plist>
