# launchd/com.kimgwonil.youtube-shorts-catchup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kimgwonil.youtube-shorts-catchup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/kimgwonil/youtube-shorts-automation/scripts/catch_up_upload_on_boot.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kimgwonil/youtube-shorts-automation</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/kimgwonil/youtube-shorts-automation/logs/catch-up-launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kimgwonil/youtube-shorts-automation/logs/catch-up-launchd.err.log</string>
</dict>
</plist>
